# supervisord.conf
[program:snorby]
command = /start.sh
autorestart = true